# eqemu-server/files/etc/supervisor/conf.d/loginserver.conf
[program:loginserver]
command=/root/server/loginserver
autorestart=true
