# eqemu-server/files/etc/supervisor/conf.d/ucs.conf
[program:ucs]
command=/root/server/ucs
autorestart=true
