# eqemu-server/files/etc/supervisor/conf.d/world.conf
[program:world]
command=/root/server/world
autorestart=true
