# eqemu-server/files/etc/supervisor/conf.d/zone.conf
[program:zone]
command=/root/server/zone
autorestart=true
process_name=%(program_name)s%(process_num)s
numprocs=2
